# caddy-mtls/rootfs/etc/s6-overlay/s6-rc.d/cert-monitor/finish
#!/command/execlineb -S0
foreground { s6-sleep 5 }
exit 0
